# projects/proj2_journal_admin/debian/simply-journal-admin.service
[Unit]
Description=simply_journal_admin recent-entries reporter
Documentation=https://github.com/ValentinTwin1206/modern-python-devops-egineering
After=systemd-journald.service
Requires=systemd-journald.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/simply-journal-admin --since-minutes 60 --format json
StandardOutput=journal
StandardError=journal
