# projects/proj2_journal_admin/debian/simply-journal-admin.service
[Unit]
Description=Simply Journal Admin journal reporter
Documentation=https://github.com/ValentinTwin1206/modern-python-devops-egineering
After=systemd-journald.service
Requires=systemd-journald.service

[Service]
# Long-running follow mode: continuously tail the journal as JSON.
Type=simple
ExecStart=/opt/simply-journal-admin/venv/bin/simply-journal-admin --follow --format json
Restart=on-failure
RestartSec=5

# --- Sandboxing / hardening -------------------------------------------------
# Run as a transient unprivileged user that can still read the system journal.
DynamicUser=yes
SupplementaryGroups=systemd-journal
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

